# docs/gunicorn-sinfoniaV2.socket
[Unit]
Description=gunicorn sinfoniaV2 socket

[Socket]
ListenStream=/run/gunicorn-sinfoniaV2.sock

[Install]
WantedBy=sockets.target
